# systemd/opencode-session-manager.service
[Unit]
Description=OpenCode Session Manager
After=network.target

[Service]
Type=simple
ExecStart=%h/.opencode-session-manager/bin/opencode-sessions serve
Restart=on-failure
Environment="PATH=__USER_PATH__"

[Install]
WantedBy=default.target
